# etc/cron.d/valkey-demo
# The first element of the path is the directory where the pluto binary is located
PATH=/usr/local/bin:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

15,35,55 * * * * valkey-planets /usr/bin/flock -n /home/valkey-planets/valkey-demo/run/cronlock /home/valkey-planets/valkey-demo/sh/pluto_fetch_n_merge.sh wait
